# jrasp-daemon/bin/jrasp-daemon.service
[Unit]
Description=jrasp-daemon service

[Service]
Type=simple
WorkingDirectory=/usr/local/jrasp/bin
ExecStart=/usr/local/jrasp/bin/jrasp-daemon
ExecStop=/bin/kill $MAINPID
Restart=always
StandardOutput=null
StandardError=null

[Install]
WantedBy=multi-user.target
